# auctioneer.service
# Copy to /etc/systemd/system/ and edit paths if your install is elsewhere.
# Then: sudo systemctl daemon-reload && sudo systemctl enable auctioneer && sudo systemctl start auctioneer

[Unit]
Description=Auctioneer Discord bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=muthu
Group=muthu
WorkingDirectory=/home/muthu/auctioneer
EnvironmentFile=/home/muthu/auctioneer/.env
ExecStart=/home/muthu/auctioneer/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=auctioneer

[Install]
WantedBy=multi-user.target
